# packages/examples/pug/crontab/deleteLogs.sh
#! /bin/bash

source "/home/q/www/${APP_CODE}/webapps/node_modules/nice-node/scripts/delete_logs.sh"
